# src/SambaScan/cron.d
# Scan for samba servers every 10 minutes and report them to the router

*/10 * * * *	root	/usr/pluto/bin/nbtsharescan.sh 

