# external/alsa-driver-1.0.10/debian/alsa-modules.postrm
#!/bin/sh

set -e

#DEBHELPER#

if [ "$1" = "remove" ]; then
	depmod -a -F /boot/System.map-_KVERS_ _KVERS_ || :
	#if [ "$(uname -r)" = "_KVERS_" ]; then
	#	invoke-rc.d alsa force-unload
	#fi
fi

