# external/hal-0.5.7/debian/hal.postrm
#!/bin/sh

set -e

# Do not restart dbus on upgrades, only if we remove the package.
if [ -x /etc/init.d/dbus ]; then
  invoke-rc.d dbus force-reload || true
fi

#DEBHELPER#

if [ "$1" = "purge" ] ; then
        deluser --quiet --system hal > /dev/null || true
        delgroup --quiet --system hal > /dev/null || true
fi

exit 0
