# src/Update/init.d/LMCEUpdate
#!/bin/bash
### BEGIN INIT INFO
# Provides:          lmceupdate
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:
# Default-Start:     2 3 4 5
# Default-Stop:      1 6
# Short-Description: LMCEUpdate
# Description:       This script updates the LMCE MD
### END INIT INFO #

# Source fuction library
. /lib/lsb/init-functions

start() {
	/usr/pluto/bin/LMCEUpdate_Apply.sh
}

### main logic ###
case "$1" in
  start)
        log_daemon_msg "Starting LMCEUpdate" "LMCEUpdate_Apply.sh"
	$1
        retval=0
        log_end_msg $retval
        ;;
  stop)
        : # noop
        retval=0
        ;;
  status)
        status_of_proc LMCEUpdate_Apply.sh LMCEUpdate
        retval=$?
        ;;
  restart|reload|force-reload)
        : # noop
        retval=0
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|reload|status}"
        retval=1
esac
exit $retval
