# external/mythtv/mythtv/debian/mythtv-backend.cron.daily
#!/bin/sh

test -x /usr/bin/mythfilldatabase || exit 0

su - mythtv -c "mythfilldatabase --quiet"
