# external/mythtv/mythtv/debian/mythtv-backend.cron.daily
#!/bin/sh

su mythtv -c "mythfilldatabase --quiet"
