# src/Update/lmce-update-system.cron.d
13 *    * * *   root	/usr/pluto/bin/LMCEUpdate_DownloadXml.sh
