# src/AVWizard/init.d/check_avwizard
#!/bin/bash
### BEGIN INIT INFO 
# Provides:          check_avwizard
# Required-Start:    
# Required-Stop:     
# Should-Start:      
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: CheckAVWizard
# Description:       This script starts the AV Wizard IF the SHIFT-KEY is pressed during startup.
# 
### END INIT INFO #

if [ "$1" == "status" ] ; then
	exit 0
fi

if [[ -x /usr/pluto/bin/WatchGyroRemote ]] ;then
	/usr/pluto/bin/WatchGyroRemote &>/dev/null </dev/null &
fi

beep -l 50 -f 2200 -n -l 50 -f 2000 -n -l 50 -f 1800 -n -l 50 -f 1600 -n -l 50 -f 1400 -n -l 50 -f 1200 -n -l 50 -f 1000 &
sleep 2

if [[ -x /usr/pluto/bin/shift_state ]] ;then
	if [[ "$( /usr/pluto/bin/shift_state )" == "1" ]] ;then
		beep -l 500
		. /usr/pluto/bin/Config_Ops.sh 
		ConfSet "AVWizardOverride" "1"
	fi
fi

# workaround for the MSI K9NGM3 motherboard
if lspci -nnv -s 00:07.0 -d 10de:055c|grep -q 1462:7349; then
	set -x
	exec 4>&1 5>&2
	exec 1>>/var/log/pluto/intel-msi-hack.log 2>&1
	date
	modprobe -r snd_hda_intel || :
	modprobe snd_hda_intel model=6stack-dell || :
	set +x
	exec 1>&4 2>&5 4>&- 5>&-
fi

