# src/Tribunexml/cron.d/TribuneTV
10 0 * * */7		root		/usr/pluto/bin/headendtobase
20 0 * * *		root		/usr/pluto/bin/Tribunexml
