# external/alsa-utils-1.0.10/debian/udev.rules
KERNEL=="controlC[0-7]", ACTION=="add", RUN+="/lib/alsa-utils/udev"
