# ubuntu/asterisk/add/debian-ubuntu/postinst
#!/bin/bash
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
#

case "$1" in
    configure)
		#enable firewal for SIP and IAX2
		echo "insert into Firewall (Protocol,SourcePort,SourcePortEnd,DestinationPort,RuleType) values ('udp',4569,0,0,'core_input');" | /usr/bin/mysql pluto_main || /bin/true
		echo "insert into Firewall (Protocol,SourcePort,SourcePortEnd,DestinationPort,RuleType) values ('udp',5060,0,0,'core_input');" | /usr/bin/mysql pluto_main || /bin/true
		echo "insert into Firewall (Protocol,SourcePort,SourcePortEnd,DestinationPort,RuleType) values ('udp',2000,0,0,'core_input');" | /usr/bin/mysql pluto_main || /bin/true
		echo "insert into Firewall (Protocol,SourcePort,SourcePortEnd,DestinationPort,RuleType) values ('tcp',2000,0,0,'core_input');" | /usr/bin/mysql pluto_main || /bin/true
		#run setup firewall
		/usr/pluto/bin/Network_Firewall.sh

		#disable oss module to speak in the house
		sed -r -i "s/.*(noload [=][>] chan_oss.so)/\1/" /etc/asterisk/modules.conf
		sed -r -i "s/.*(noload [=][>] chan_alsa.so)/\1/" /etc/asterisk/modules.conf
		sed -r -i "s/^load[ ][=][>][ ]chan_modem\.so//" /etc/asterisk/modules.conf

		if getent passwd asterisk >/dev/null; then
			# fix home directory if needed
			if [[ "$(getent passwd asterisk|awk -F: '{print $6}')" == "/var/run/asterisk" ]];then
				usermod -d /var/lib/asterisk asterisk || /bin/true
			fi
		else
			# add user if it doesn't exist
			adduser --system --group --home /var/lib/asterisk --gecos "Asterisk PBX daemon" asterisk || /bin/true
			adduser asterisk audio || /bin/true
		fi

		echo 'CREATE TABLE IF NOT EXISTS devices (id VARCHAR( 20 ) NOT NULL , tech VARCHAR( 10 ) NOT NULL , dial VARCHAR( 50 ) NOT NULL , devicetype VARCHAR( 5 ) NOT NULL , user VARCHAR( 50 ) , description VARCHAR( 50 ));' | /usr/bin/mysql asterisk || /bin/true
		echo 'CREATE TABLE IF NOT EXISTS users (extension VARCHAR( 20 ) NOT NULL , password VARCHAR( 20 ) , name VARCHAR( 50 ) , voicemail VARCHAR( 50 ) , ringtimer INT(3) , noanswer VARCHAR( 100 ) , recording VARCHAR( 50 ) ,  outboundcid VARCHAR( 50 ));' | /usr/bin/mysql asterisk || /bin/true
		echo 'CREATE TABLE IF NOT EXISTS `incoming` ( `cidnum` VARCHAR( 20 ) , `extension` VARCHAR( 20 ) , `destination` VARCHAR( 50 ) , `faxexten` VARCHAR( 20 ) , `faxemail` VARCHAR( 50 ) , `answer` TINYINT( 1 ) , `wait` INT( 2 ) , `privacyman` TINYINT( 1 ) );' | /usr/bin/mysql asterisk || /bin/true
		echo 'CREATE TABLE IF NOT EXISTS `Backup` (`Name` varchar(50) default NULL,`Voicemail` varchar(50) default NULL,`Recordings` varchar(50) default NULL,`Configurations` varchar(50) default NULL,`CDR` varchar(55) default NULL,`FOP` varchar(50) default NULL,`Minutes` varchar(50) default NULL,`Hours` varchar(50) default NULL,`Days` varchar(50) default NULL,`Months` varchar(50) default NULL,`Weekdays` varchar(50) default NULL,`Command` varchar(200) default NULL,`Method` varchar(50) default NULL,`ID` int(11) NOT NULL auto_increment,PRIMARY KEY  (ID));' | /usr/bin/mysql asterisk || /bin/true
		sed -r -i "s/^autoload=yes/autoload=yes\nnoload => chan_skinny.so\nnoload => chan_capi.so\nnoload => app_capiRETRIEVE.so\nnoload => app_capiCD.so\nnoload => app_capiECT.so\nnoload => app_capiFax.so\nnoload => app_capiHOLD.so\nnoload => app_capiMCID.so\nnoload => app_capiNoES.so/" ${PKGFOLDER}/etc/asterisk/modules.conf 
		sed -r -i "s/.*(noload [=][>] chan_oss.so)/\1/" /etc/asterisk/modules.conf
		sed -r -i "s/.*(noload [=][>] chan_alsa.so)/\1/" /etc/asterisk/modules.conf
		killall asterisk safe_asterisk || /bin/true

		mkdir -p /var/lib/asterisk/sounds/voicemail
		chmod -R ug+rx /var/lib/asterisk/sounds/voicemail
		
		### disable modem 
		echo "Disable modem"
		sed -r -i "s/^load[ ][=][>][ ]chan_modem\.so//" /etc/asterisk/modules.conf
		
		chmod -R ug+rw /etc/asterisk/
		chmod 777 /var /var/lib/ /var/lib/asterisk/ /var/lib/asterisk/sounds/
		chmod -R ug+rw /var/lib/asterisk/sounds/

		if [[ "$(pidof asterisk)" != "" ]] ;then
			invoke-rc.d asterisk restart || /bin/true
		fi
		
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)

    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
