# src/Embedded_Boot/etc/init.d/check_internet
#!/bin/sh /etc/rc.common

START=96

boot()
{
	if [[ ! -f /etc/pluto/fresh_install ]]; then
		return 0
	fi
	echo 0 >/proc/diag/led/power
	/usr/bin/NetLED.sh >/dev/null 2>&1 </dev/null &
}
