# src/Embedded_Boot/etc/udhcpc.user
#!/bin/sh /etc/rc.common

if [[ ! -f /etc/pluto/fresh_install ]]; then
	/usr/bin/Report_AP.sh &
	return 0
fi

WANip="$ip"

IPtype=routable

d1=$(echo "$WANip" | cut -d. -f1)
d2=$(echo "$WANip" | cut -d. -f2)

if [[ "$d1" == 10 ]]; then
	IPtype=local
elif [[ "$d1" == 192 ]] && [[ "$d2" == 168 ]]; then
	IPtype=local
elif [[ "$d1" == 172 ]] && [[ "$d2" -ge 16 ]] && [[ "$d2" -le 31 ]]; then
	IPtype=local
fi

Cfg=$(uci show dhcp|grep interface=lan|cut -d. -f2)
CfgMasq=$(uci show dhcp|grep =dnsmasq|cut -d. -f2|cut -d= -f1)
uci set dhcp.$Cfg.force=1
uci set dhcp.$Cfg.ignore=0
case "$IPtype" in
	local)
		uci set dhcp.$Cfg.dynamicdhcp=0
		uci set dhcp.$CfgMasq.authoritative=0
		uci set network.eth0.vlan0="0 1 2 3 4 5*"
		uci del network.eth0.vlan1
		uci set network.wan.ifname=br-lan:0
		;;
	routable)
		uci set dhcp.$Cfg.dynamicdhcp=1
		uci set dhcp.$CfgMasq.authoritative=1
		
		Model=$(cat /proc/diag/model)
		case "$Model" in
			"ASUS WL-500g Premium")
				uci set network.eth0.vlan0="1 2 3 4 5*"
				uci set network.eth0.vlan1="0 5"
				;;
			"ASUS WL-500g Premium V2")
				uci set network.eth0.vlan0="0 1 2 3 5*"
				uci set network.eth0.vlan1="4 5"
				;;
		esac

		uci set network.wan.ifname=eth0.1
		;;
esac
uci commit

rm -f /etc/pluto/fresh_install

/etc/init.d/firewall restart
/etc/init.d/dnsmasq restart
/etc/init.d/network restart

/usr/bin/Report_AP.sh &
