# ubuntu/lirc-0.8.2+lmce/debian/lirc.preinst
#!/bin/sh

set -e

#DEBHELPER#

case "$1" in
    upgrade)
	if dpkg --compare-versions "${2}" le-nl "0.8.2-0ubuntu6"
	then update-rc.d -f lirc remove > /dev/null
	fi
        ;;
esac


