# src/Add_Software/cron.d/Add_Software
0 0 * * */7		root		/usr/pluto/bin/getxmls
