# src/Add_Software/cron.d/Add_Software
0 0 * * */7		root		test -f /etc/diskless.conf || bash -c '/usr/pluto/bin/getxmls &>/dev/null'
